# ftg-master2.service
[Unit]
Description=Follow The Goat - Master2 (Trading Logic)
After=network.target postgresql.service ftg-master.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/follow_the_goat
Environment="PATH=/root/follow_the_goat/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/follow_the_goat/venv/bin/python scheduler/master2.py
Restart=always
RestartSec=10
StandardOutput=append:/root/follow_the_goat/logs/scheduler2.log
StandardError=append:/root/follow_the_goat/logs/scheduler2_errors.log

[Install]
WantedBy=multi-user.target
